# themes/marvell/root/etc/uci-defaults/luci-theme-marvell
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Marvell=/luci-static/marvell
	commit luci
EOF
exit 0
